# deploy/halloween-party.service
[Unit]
Description=Halloween Karaoke Costume Contest
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=halloween
Group=halloween
WorkingDirectory=/opt/halloween/current
Environment=PYTHONUNBUFFERED=1
Environment=APP_ENV=production
Environment=VAULT_ADDR=http://172.31.118.0:8200
Environment=VAULT_AUTH_METHOD=aws
Environment=VAULT_AWS_AUTH_ROLE=goodvines-api
Environment=AWS_REGION=us-east-1
Environment=HALLOWEEN_APP_PORT=8081
ExecStart=/opt/halloween/current/deploy/start_halloween.sh
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
